# source/04.servidor/07.tareas/files/mpvcleaner.timer
[Unit]
Descripcion=Ejecuta al abrir la primera sesión de usuario mpvcleaner.sh

[Timer]
OnStartupSec=5s

[Install]
WantedBy=timers.target
